# deploy/led-ticker.service
[Unit]
Description=LED Score Ticker Display
After=score-fetcher.service
Wants=score-fetcher.service

[Service]
Type=simple
# Must run as root for GPIO access
User=root
WorkingDirectory=/home/pi/sports-score-ticker
ExecStart=taskset -c 3 /home/pi/sports-score-ticker/display/build/led_ticker --led-no-hardware-pulse --led-slowdown-gpio=4 --led-pwm-lsb-nanoseconds=200 --led-pwm-bits=9
Restart=always
RestartSec=5
Environment=TICKER_DATA_FILE=/tmp/scores.json

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=led-ticker

# Security hardening (still needs GPIO)
ProtectHome=read-only
ProtectSystem=strict
ReadWritePaths=/tmp
PrivateTmp=false

[Install]
WantedBy=multi-user.target
