# deploy/led-ticker.service
[Unit]
Description=LED Score Ticker Display
After=score-fetcher.service
Wants=score-fetcher.service

[Service]
Type=simple
# Must run as root for GPIO access
User=root
WorkingDirectory=/home/pi/sports-score-ticker
ExecStart=/home/pi/sports-score-ticker/display/build/led_ticker --led-rows=32 --led-cols=64 --led-chain=5 --led-parallel=2 --led-hardware-mapping=regular --led-slowdown-gpio=2 --led-brightness=80 --led-pwm-lsb-nanoseconds=130
Restart=always
RestartSec=5
Environment=TICKER_DATA_FILE=/tmp/scores.json

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=led-ticker

# Security hardening (still needs GPIO)
ProtectHome=read-only
ProtectSystem=strict
ReadWritePaths=/tmp
PrivateTmp=false

[Install]
WantedBy=multi-user.target
